# files/profile.d/aliases.sh
#!/bin/bash

alias psg="ps aux | grep"
alias pgg="pgrep -lf"
alias k9="kill -9"
alias rm="rm -i"
alias cp="cp -i"
