# files/profile.d/aliases.sh
#!/bin/bash

alias ls="ls -G"
alias l="ls -h"
alias ll="ls -lh --group-directories-first"
alias la="ls -alh --group-directories-first"
alias psg="ps aux | grep"
alias pgg="pgrep -lf"
alias k9="kill -9"
alias rm="rm -i"
alias cp="cp -i"
