# files/profile.d/aws.sh
#!/bin/bash

# AWSCLI
complete -C aws_completer aws
